# SysVolumeServer/rc.local
#!/bin/bash

/SysVolumeServer/nfs_ServiceStart.sh

exit 0
