# DCloudSwift/misc/BootScripts/rc.local
#!/bin/sh
sleep 10
python /etc/rc.py

exit 0
