# DCloudSwift/misc/BootScripts/rc.local
sleep 10
sh /etc/DCloud/ServerRegister/autoRun.sh
python /etc/rc.py

exit 0
