# DCloudGateway/gateway_scripts/set_proxy.conf
description     "upstart script for setting squid3 status"
author          "Weitang Hung <weitang.hung@delta.com.tw"

start on (started squid3)

exec /etc/delta/set_proxy.py
